# packages/lattica-polymorph/lattica_polymorph-0.1.2.tar.gz/lattica_polymorph-0.1.2/monitor/exporters/node_exporter/node_exporter.service
[Unit]
Description=Prometheus Node Exporter
After=network.target

[Service]
User=william
ExecStart=/usr/local/bin/node_exporter
Restart=always

[Install]
WantedBy=multi-user.target
