# packages/naneos-devices/naneos_devices-1.1.16.tar.gz/naneos_devices-1.1.16/installers/rp-naneos-uploader/naneos_uploader.service
[Unit]
Description=Naneos Uploader Service
After=network.target

[Service]
ExecStart={{APP_DIR}}/.venv/bin/python {{APP_DIR}}/uploader-script.py
WorkingDirectory={{APP_DIR}}
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=5
User={{USER_NAME}}

[Install]
WantedBy=multi-user.target