# packages/autonomous-coding/autonomous_coding-2.0.3.tar.gz/autonomous_coding-2.0.3/docs/com.autonomous.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
Orchestrator LaunchAgent Plist (macOS)
======================================

Installation:
  1. Copy this file to ~/Library/LaunchAgents/com.autonomous.orchestrator.plist
  2. Edit the paths below to match your setup
  3. Load the agent: launchctl load ~/Library/LaunchAgents/com.autonomous.orchestrator.plist

Commands:
  - Load: launchctl load ~/Library/LaunchAgents/com.autonomous.orchestrator.plist
  - Unload: launchctl unload ~/Library/LaunchAgents/com.autonomous.orchestrator.plist
  - Start: launchctl start com.autonomous.orchestrator
  - Stop: launchctl stop com.autonomous.orchestrator
  - Check status: launchctl list | grep orchestrator
  - View logs: tail -f /tmp/orchestrator.stdout.log

Notes:
  - LaunchAgents run in user context (not root)
  - For system-wide service, use /Library/LaunchDaemons instead
  - Ensure Python and all dependencies are available
-->
<plist version="1.0">
<dict>
    <!-- Service identifier (must match filename) -->
    <key>Label</key>
    <string>com.autonomous.orchestrator</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/path/to/autonomous-coding</string>

    <!-- Command to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>orchestrator.py</string>
        <string>--project-dir</string>
        <string>/Users/YOUR_USERNAME/path/to/project</string>
    </array>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>AGENT_TYPE</key>
        <string>ORCHESTRATOR</string>
        <!-- Add API keys here -->
        <!-- <key>ANTHROPIC_API_KEY_1</key>
        <string>your-key-here</string> -->
    </dict>

    <!-- Run at login (user agent) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep running -->
    <key>KeepAlive</key>
    <dict>
        <!-- Restart if crashed -->
        <key>SuccessfulExit</key>
        <false/>
        <!-- Don't restart if manually stopped -->
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Throttle restarts (wait 10 seconds) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/orchestrator.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/orchestrator.stderr.log</string>

    <!-- Nice level (lower priority) -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Soft termination before SIGKILL -->
    <key>ExitTimeOut</key>
    <integer>30</integer>
</dict>
</plist>
