# packages/pglift/pglift-2.6.1.tar.gz/pglift-2.6.1/src/pglift/systemd/pglift-backup@.timer
[Unit]
Description=Backup %i PostgreSQL database instance

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
FixedRandomDelay=true
Persistent=true

[Install]
WantedBy=timers.target
