# packages/pglift/pglift-2.6.1.tar.gz/pglift-2.6.1/src/pglift/systemd/pglift-backup@.service
[Unit]
Description=Backup %i PostgreSQL database instance

[Service]
Type=oneshot
{executeas}

{environment}
ExecStart={execpath} instance backup %I
