# packages/polarion-mcp/polarion_mcp-2.0.0.tar.gz/polarion_mcp-2.0.0/polarion-mcp.service
[Unit]
Description=Polarion MCP Server
After=network.target

[Service]
Type=simple
User=support
WorkingDirectory=/home/support/polarion-mcp-server
EnvironmentFile=/home/support/polarion-mcp-server/.env
ExecStart=/home/support/.local/bin/polarion-mcp
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

