# packages/pglift/pglift-2.6.1-py3-none-any.whl/pglift/systemd/pglift-temboard_agent@.service
[Unit]
Description=temBoard agent for %i
After=network.target pglift-postgresql@%i.service

[Service]
Type=simple
{executeas}
ExecStart={execpath} -c {configpath}

[Install]
WantedBy=default.target
