# packages/pglift/pglift-2.6.1-py3-none-any.whl/pglift/systemd/pglift-postgres_exporter@.service
[Unit]
Description=Prometheus exporter %i database server metrics
Documentation=https://github.com/prometheus-community/postgres_exporter

[Service]
{executeas}
EnvironmentFile=-{configpath}
ExecStart={execpath} $POSTGRES_EXPORTER_OPTS
Restart=on-failure

[Install]
WantedBy=default.target
