# packages/lattica-polymorph/lattica_polymorph-0.1.2.tar.gz/lattica_polymorph-0.1.2/monitor/exporters/scraper/exporter.service
[Unit]
Description="Scraper Prometheus exporter"
After=network.target

[Service]
User=william
ExecStart=/usr/bin/python3 /opt/polymorph_exporter/polymorph_exporter.py
Restart=always

[Install]
WantedBy=multi-user.target
