# packages/lattica-polymorph/lattica_polymorph-0.1.2.tar.gz/lattica_polymorph-0.1.2/monitor/exporters/scraper/scraper.service
[Unit]
Description="Scraper runner"
After=network.target

[Service]
User=william
ExecStart=/usr/local/bin/run_scrape.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
