# packages/mirror-dedupe/mirror_dedupe-0.2.10.tar.gz/mirror_dedupe-0.2.10/systemd/mirror-dedupe.service
[Unit]
Description=Mirror Deduplication Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/mirror-dedupe --config /etc/mirror-dedupe
StandardOutput=journal
StandardError=journal
