# packages/lattica-polymorph/lattica_polymorph-0.1.2.tar.gz/lattica_polymorph-0.1.2/monitor/exporters/tests/tests_exporter.service
[Unit]
Description="Polymorph Tests Prometheus Exporter"
After=network.target

[Service]
User=ubuntu
ExecStart=/usr/bin/python3 /opt/export/tests_exporter.py
Restart=always

[Install]
WantedBy=multi-user.target
