# packages/jdcat/jdcat-1.2.0-py3-none-any.whl/jdcat-1.2.0.data/data/share/jdcat/resources/com.jdcat.proxy.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jdcat.proxy.helper</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Library/PrivilegedHelperTools/jdcat_proxy_helper.py</string>
      <string>--port</string>
      <string>17901</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>UserName</key>
    <string>root</string>

    <!-- Optional logs for troubleshooting -->
    <key>StandardOutPath</key>
    <string>/var/log/com.jdcat.proxy.helper.out.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/com.jdcat.proxy.helper.err.log</string>
  </dict>
</plist>