# packages/risclog.batou/risclog_batou-3.7.4.tar.gz/risclog.batou-3.7.4/src/risclog/batou/resources/proftpd/proftpd.service
# Managed by batou
# {{component.config_hash}}

[Unit]
Description=Proftpd

[Service]
Environment="LOCALE_ARCHIVE={{component.env['LOCALE_ARCHIVE']}}"
Environment="PATH={{component.env['PATH']}}"
Environment="TZDIR={{component.env['TZDIR']}}"
Type=forking

User=root
Group=root

ExecStart=/srv/{{component.service_user}}/.nix-profile/bin/proftpd -c {{component.config_file.path}}
